# debian/postinst
#!/bin/sh
# postinst script for langstroth
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

ucf --debconf-ok /usr/share/langstroth/settings_example.py /etc/langstroth/settings.py
ucf --debconf-ok /usr/share/langstroth/langstroth.apache.conf /etc/langstroth/apache.conf

CONF="langstroth.conf"

# Only try to add a symlink on a fresh install to respect
# changes done by the administrator
if [ "$2" = '' ]; then
    server="apache2"
    if [ -d "/etc/${server}/conf-available" ]; then
        if [ ! -e "/etc/${server}/conf-available/$CONF" ] ; then
            ln -s ../../langstroth/apache.conf "/etc/${server}/conf-available/$CONF"
        fi
    fi
    if [ -e /usr/share/apache2/apache2-maintscript-helper ] ; then
        # apache 2.4
        . /usr/share/apache2/apache2-maintscript-helper
        apache2_invoke enconf $CONF || exit $?
    fi
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

django-admin migrate --settings="langstroth.settings" --noinput
django-admin collectstatic --settings="langstroth.settings" --noinput

chown -R www-data:www-data /var/log/langstroth
chmod -R u+w /var/log/langstroth

mkdir -p /var/lib/langstroth/static/dashboard/js
chown -R www-data:www-data /var/lib/langstroth
chmod -R u+w /var/lib/langstroth/static

if [ -f /etc/init.d/apache2 ]
then
    invoke-rc.d apache2 reload
fi

exit 0
